# docker/easypanel/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
autorestart=true

[program:apache]
command=apache2-foreground
autorestart=true
